# setup/config/nginx.sv.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority=2

[eventlistener:nginx_exit]
command=/setup/kill_supervisor.py
process_name=nginx
events=PROCESS_STATE_EXITED
