# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/openvpn --config /config/openvpn/openvpn.ovpn
